# bifrost/data/systemd/bifrost-save-sync.service
[Unit]
Description=Bifrost save/state sync — local ↔ RomM
Documentation=https://github.com/rommapp/romm-bifrost
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/bifrost save-sync --apply
ExecStart=%h/.local/bin/bifrost state-sync --apply
StandardOutput=journal
StandardError=journal
RemainAfterExit=no

[Install]
WantedBy=default.target
